# product/server/codekvast-server/src/dist/etc/startscript/gentoo/codekvast-server
#!/sbin/runscript

depend() {
    need net
}

start() {
    ebegin "Starting codekvast_server"
    start-stop-daemon --start --exec /opt/codekvast-server/bin/codekvast-server \
        --background --make-pidfile --pidfile /var/run/codekvast-server.pid \
        --stdout /var/log/codekvast/codekvast-server.out --stderr /var/log/codekvast/codekvast-server.out
    eend $?
}

stop() {
    ebegin "Stopping codekvast_server"
    start-stop-daemon --stop --exec /opt/codekvast-server/bin/codekvast-server --pidfile /var/run/codekvast-server.pid
    eend $?
}
