# product/warehouse/src/dist/etc/startscript/redhat/codekvast-warehouse
#!/bin/bash
#
# codekvast-warehouse     This shell script controls the codekvast warehouse.
#
# Author:       Olle Hallin <olle.hallin@crisp.se>
#
# chkconfig:    345 97 2
#
# description:  Processes data from the Codekvast collectors
# processname:  codekvast-warehouse
#

# source function library
. /etc/rc.d/init.d/functions

BINDIR=/opt/codekvast-warehouse/bin
LOGDIR=/var/log/codekvast
PROG=codekvast-warehouse
LOCK=/var/lock/subsys/$PROG
PIDFILE=/var/run/$PROG.pid

RETVAL=0

start() {
  if [ ! -f $LOCK ]; then
    echo -n "Starting $PROG..."
    mkdir -p $LOGDIR
    nohup $BINDIR/$PROG < /dev/null > $LOGDIR/$PROG.out 2>&1 &
    echo $! > $PIDFILE
    RETVAL=$?
    [ $RETVAL -eq 0 ] && touch $LOCK
  else
    echo -n "$PROG is already running."
  fi
  echo
  return $RETVAL
}

stop() {
  echo -n "Stopping $PROG..."
  killproc -p "$PIDFILE" $PROG || failure
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && rm -f $LOCK
  return $RETVAL
}

restart() {
  stop
  start
}

reload() {
    restart
}

case "$1" in
  start)
      start
      ;;
  stop)
      stop
      ;;
  restart)
      restart
      ;;
  condrestart)
      [ -e $LOCK ] && restart || :
      ;;
  status)
      status -p "$PIDFILE" $PROG
      RETVAL="$?"
      ;;
  *)
  echo $"Usage: $0 {start|stop|status|restart|condrestart}"
  exit 1
esac

exit $RETVAL
