# product/agent/codekvast-agent/src/dist/etc/startscript/gentoo/codekvast-agent
#!/sbin/runscript

depend() {
    need net
}

start() {
    ebegin "Starting codekvast_agent"
    start-stop-daemon --start --exec /opt/codekvast-agent/bin/codekvast-agent \
        --background --make-pidfile --pidfile /var/run/codekvast-agent.pid \
        --stdout /var/log/codekvast/codekvast-agent.out --stderr /var/log/codekvast/codekvast-agent.out
    eend $?
}

stop() {
    ebegin "Stopping codekvast_agent"
    start-stop-daemon --stop --exec /opt/codekvast-agent/bin/codekvast-agent --pidfile /var/run/codekvast-agent.pid
    eend $?
}
